# rootfs/etc/systemd/system/zero-webapp.service
[Unit]
Description=Zero Web Application
After=network-online.target zero-wifi-portal.service
Wants=network-online.target
ConditionPathExists=/etc/zero/.wifi_configured

[Service]
Type=simple
User=root
WorkingDirectory=/opt/zero/web
Environment="PATH=/opt/zero/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/etc/zero/secrets.env
ExecStart=/opt/zero/.venv/bin/python app.py
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/zero /var/log

[Install]
WantedBy=multi-user.target
