# rootfs/etc/systemd/system/zero-updater.service
[Unit]
Description=Zero OTA Update Check
After=network-online.target
Wants=network-online.target
# Don't run during initial provisioning
ConditionPathExists=!/etc/zero/.provisioning

[Service]
Type=oneshot
ExecStart=/opt/zero/scripts/update.sh update
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
