# rootfs/etc/systemd/system/zero-wifi-portal.service
[Unit]
Description=Zero WiFi Configuration Portal
After=network.target NetworkManager.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/zero/wifi-portal
Environment="PATH=/opt/zero/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/etc/zero/secrets.env
ExecStart=/opt/zero/.venv/bin/python app.py
Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/zero /var/log

[Install]
WantedBy=multi-user.target
